# deploy/viewport.service
[Unit]
Description=viewport — headless RTSP multi-stream display
After=network-online.target
Wants=network-online.target

[Service]
# Replace with the name of the unprivileged service user you created.
# If you followed the README, this is 'viewport'.
User=viewport
# Primary group — grants access to /dev/dri/* (DRM/KMS output).
Group=video
# Some systems assign /dev/dri/renderD128 to the 'render' group instead of
# (or in addition to) 'video'.  Comment out if not required for your system.
SupplementaryGroups=render
Environment=PYTHONUNBUFFERED=1

# Adjust the path to your config file as needed
ExecStart=/usr/bin/python3 /opt/viewport/main.py /etc/viewport/config.yaml

# Restart on any failure (including stream errors), with a short delay
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

